# etc/systemd/system/qBitPuller.service
[Unit]
Description=Seedbox puller (qBittorrent API + rclone)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=PYTHONIOENCODING=utf-8
ExecStart=/usr/bin/python3 /usr/local/bin/qBitPuller.py
